# conf/pind.service
[Unit]
Description=Pind - pin daemon
After=network.target
Wants=network.target
Requires=network.target

[Service]
Type=simple
User=root
Group=root
LimitNOFILE=500000

ExecStart=/usr/bin/pind --config /etc/pind/pind.yml --service
ExecReload=/usr/bin/kill -HUP $MAINPID

Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
